# deploy/systemd/pm-safeline-deploy.service
[Unit]
Description=PM 세이프라인 자동 배포 (git pull + 빌드 + 재시작)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=deploy
WorkingDirectory=/opt/pm-safeline/app
# 빌드 시 gradle 이 java 를 찾을 수 있도록.
Environment=JAVA_HOME=/usr/lib/jvm/java-21-openjdk
Environment=PATH=/usr/lib/jvm/java-21-openjdk/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/opt/pm-safeline/app/deploy/deploy.sh
