# deploy/systemd/pm-safeline.service
[Unit]
Description=PM 세이프라인 안전경로 서버 (Ktor/GraphHopper)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/opt/pm-safeline/app

# JDK 21 필요. 배포판 실행 스크립트가 java 를 PATH 에서 찾는다. 경로는 서버에 맞게 수정.
Environment=JAVA_HOME=/usr/lib/jvm/java-21-openjdk
Environment=PATH=/usr/lib/jvm/java-21-openjdk/bin:/usr/local/bin:/usr/bin:/bin
Environment=JAVA_OPTS=-Xmx1g

# 앱은 WorkingDirectory 의 .env 를 자동 로드한다(Dotenv). PM_OSM_FILE 등은 .env 에 두면 됨.
# (원한다면 여기서 직접 지정도 가능)
# Environment=PM_OSM_FILE=/opt/pm-safeline/app/data/osm/daejeon.osm

# installDist 산출물의 실행 스크립트. rootProject.name=safety 이므로 launcher 는 'safety'.
ExecStart=/opt/pm-safeline/app/build/install/safety/bin/safety

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
